# ansible/roles/backend/templates/flask-app.service.j2
[Unit]
Description=Flask Application Service
After=network.target

[Service]
User={{ app_user }}
WorkingDirectory={{ app_directory }}
ExecStart=/usr/bin/python3 app.py
Restart=always

[Install]
WantedBy=multi-user.target